# kkk.conf
[program:kkk]
command=python3 /home/curator2/kkk/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/kkk/test.err.log
stdout_logfile=/home/curator2/kkk/test.out.log